# tests_source/intel-driven/2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11/2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11_detections.kql
// ============================================================================
// DEFENSE GUIDANCE: KQL Detection Queries
// ============================================================================
// Test ID: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
// Test Name: SafePay UAC Bypass & Defense Evasion
// MITRE ATT&CK: T1548.002, T1562.001, T1547.001
// Created: 2025-12-07
// Author: F0RT1KA Defense Guidance Builder
// ============================================================================
//
// TECHNIQUE-FOCUSED DETECTION PRINCIPLE:
// These detections target the underlying UAC bypass, registry persistence,
// and defense evasion techniques - NOT just the F0RT1KA testing framework.
// They will catch real-world attackers using similar techniques.
//
// ============================================================================


// ============================================================================
// QUERY 1: CMSTPLUA COM Object UAC Bypass Detection
// Purpose: Detect UAC bypass via CMSTPLUA COM interface abuse
// Confidence: HIGH
// False Positives: Rare - legitimate CMSTPLUA usage is uncommon
// MITRE ATT&CK: T1548.002
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where ProcessCommandLine has_any (
    "3E5FC7F9-9A51-4367-9063-A120244FBEC7",  // CMSTPLUA GUID
    "CMSTPLUA",
    "GetTypeFromCLSID",
    "CreateInstance"
)
| where FileName in~ ("powershell.exe", "pwsh.exe", "cmd.exe", "wscript.exe", "cscript.exe")
| project TimeGenerated, DeviceName, AccountName, FileName, ProcessCommandLine,
    InitiatingProcessFileName, InitiatingProcessCommandLine
| extend
    Severity = "Critical",
    ThreatType = "UAC Bypass - CMSTPLUA COM Object Abuse",
    MitreAttack = "T1548.002",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Confidence = "High"


// ============================================================================
// QUERY 2: PowerShell Script Block Logging - CMSTPLUA Pattern
// Purpose: Detect CMSTPLUA bypass attempt via script block logs
// Confidence: HIGH
// False Positives: Very rare - this is a specific attack pattern
// Requires: PowerShell Script Block Logging enabled (Event ID 4104)
// ============================================================================

Event
| where TimeGenerated > ago(1h)
| where Source == "Microsoft-Windows-PowerShell"
| where EventID == 4104
| where RenderedDescription has_any (
    "3E5FC7F9-9A51-4367-9063-A120244FBEC7",
    "CMSTPLUA",
    "System.Activator",
    "GetTypeFromCLSID"
)
| extend ScriptBlock = extract(@"ScriptBlockText=(.+)", 1, RenderedDescription)
| project TimeGenerated, Computer, ScriptBlock, EventID
| extend
    Severity = "Critical",
    ThreatType = "UAC Bypass Script Detected",
    MitreAttack = "T1548.002",
    Confidence = "High"


// ============================================================================
// QUERY 3: PowerShell Execution Policy Bypass from c:\F0
// Purpose: Detect suspicious PowerShell execution with bypass from test directory
// Confidence: HIGH
// False Positives: Authorized security testing
// MITRE ATT&CK: T1059.001
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where FileName in~ ("powershell.exe", "pwsh.exe")
| where ProcessCommandLine has_all ("-ExecutionPolicy", "Bypass")
| where ProcessCommandLine has_any ("c:\\F0", "c:/F0", "safepay", "uac_bypass")
| project TimeGenerated, DeviceName, AccountName, FileName,
    ProcessCommandLine, InitiatingProcessFileName, FolderPath
| extend
    Severity = "High",
    ThreatType = "Suspicious PowerShell Execution - F0RT1KA Test or Similar",
    MitreAttack = "T1059.001",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Confidence = "High"


// ============================================================================
// QUERY 4: Registry Run Key Persistence - SafePay Indicator
// Purpose: Detect SafePay-specific registry persistence value
// Confidence: HIGH
// False Positives: None - this is a specific IoC
// MITRE ATT&CK: T1547.001
// ============================================================================

DeviceRegistryEvents
| where TimeGenerated > ago(1h)
| where RegistryKey has "CurrentVersion\\Run"
| where ActionType == "RegistryValueSet"
| where RegistryValueName =~ "SafePayService"
    or RegistryValueData has "6F22-C16F-0C71-688A"
| project TimeGenerated, DeviceName, RegistryKey, RegistryValueName,
    RegistryValueData, InitiatingProcessFileName, InitiatingProcessCommandLine,
    AccountName
| extend
    Severity = "Critical",
    ThreatType = "SafePay Malware Persistence Indicator",
    MitreAttack = "T1547.001",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Confidence = "High"


// ============================================================================
// QUERY 5: Generic Registry Run Key Persistence Detection
// Purpose: Detect suspicious Run key modifications from non-standard processes
// Confidence: MEDIUM
// False Positives: Legitimate software installations, startup apps
// MITRE ATT&CK: T1547.001
// ============================================================================

DeviceRegistryEvents
| where TimeGenerated > ago(1h)
| where RegistryKey has_any (
    "\\CurrentVersion\\Run",
    "\\CurrentVersion\\RunOnce",
    "\\CurrentVersion\\RunServices"
)
| where ActionType == "RegistryValueSet"
| where InitiatingProcessFileName !in~ (
    "explorer.exe",
    "msiexec.exe",
    "setup.exe",
    "installer.exe",
    "MpCmdRun.exe",
    "MsMpEng.exe",
    "svchost.exe",
    "services.exe"
)
| where InitiatingProcessFolderPath !startswith "C:\\Windows\\System32"
| where InitiatingProcessFolderPath !startswith "C:\\Program Files"
| project TimeGenerated, DeviceName, RegistryKey, RegistryValueName,
    RegistryValueData, InitiatingProcessFileName, InitiatingProcessFolderPath,
    InitiatingProcessCommandLine, AccountName
| extend
    Severity = case(
        InitiatingProcessFolderPath has "c:\\F0", "High",
        InitiatingProcessFolderPath has "Temp", "High",
        InitiatingProcessFolderPath has "AppData", "Medium",
        "Low"
    ),
    ThreatType = "Suspicious Registry Run Key Modification",
    MitreAttack = "T1547.001",
    Confidence = "Medium"


// ============================================================================
// QUERY 6: Windows Defender Tampering via PowerShell
// Purpose: Detect attempts to disable/modify Windows Defender
// Confidence: HIGH
// False Positives: IT administrators making legitimate changes
// MITRE ATT&CK: T1562.001
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where FileName in~ ("powershell.exe", "pwsh.exe")
| where ProcessCommandLine has_any (
    "Set-MpPreference",
    "DisableRealtimeMonitoring",
    "DisableBehaviorMonitoring",
    "DisableIOAVProtection",
    "DisableBlockAtFirstSeen",
    "DisableIntrusionPreventionSystem",
    "Windows Defender",
    "SecurityHealth"
)
| where ProcessCommandLine has_any (
    "$true",
    "True",
    "-Disable",
    "Stop-Service",
    "Set-Service"
)
| project TimeGenerated, DeviceName, AccountName, FileName,
    ProcessCommandLine, InitiatingProcessFileName
| extend
    Severity = "Critical",
    ThreatType = "Windows Defender Tampering Attempt",
    MitreAttack = "T1562.001",
    Confidence = "High"


// ============================================================================
// QUERY 7: GUI Automation for Security Tool Tampering
// Purpose: Detect GUI automation attempts targeting Windows Security
// Confidence: MEDIUM
// False Positives: Legitimate automation tools, accessibility software
// MITRE ATT&CK: T1562.001
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where ProcessCommandLine has_any (
    "FindWindow",
    "ShowWindow",
    "SendMessage",
    "PostMessage",
    "UIAutomation",
    "AutoIt",
    "AutoHotkey"
)
| where ProcessCommandLine has_any (
    "Windows Security",
    "SecurityHealthHost",
    "Windows Defender",
    "Virus & threat",
    "Real-time protection"
)
| project TimeGenerated, DeviceName, AccountName, FileName,
    ProcessCommandLine, InitiatingProcessFileName, FolderPath
| extend
    Severity = "High",
    ThreatType = "GUI Automation Targeting Security Tools",
    MitreAttack = "T1562.001",
    Confidence = "Medium"


// ============================================================================
// QUERY 8: Script Drop to c:\F0 Directory
// Purpose: Detect script file creation in F0RT1KA test directory
// Confidence: HIGH
// False Positives: Authorized F0RT1KA testing
// ============================================================================

DeviceFileEvents
| where TimeGenerated > ago(1h)
| where ActionType == "FileCreated"
| where FolderPath startswith "c:\\F0" or FolderPath startswith "C:\\F0"
| where FileName endswith ".ps1"
    or FileName endswith ".bat"
    or FileName endswith ".cmd"
    or FileName endswith ".vbs"
    or FileName endswith ".js"
| project TimeGenerated, DeviceName, AccountName, FileName, FolderPath,
    FileSize = FileSize, InitiatingProcessFileName, InitiatingProcessCommandLine
| extend
    Severity = "High",
    ThreatType = "Script File Created in F0RT1KA Test Directory",
    MitreAttack = "T1105",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Confidence = "High"


// ============================================================================
// QUERY 9: UAC Bypass via Known Auto-Elevate Binaries
// Purpose: Detect UAC bypass attempts using common auto-elevate vectors
// Confidence: MEDIUM
// False Positives: Legitimate use of these tools
// MITRE ATT&CK: T1548.002
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where FileName in~ (
    "eventvwr.exe",
    "fodhelper.exe",
    "sdclt.exe",
    "computerdefaults.exe",
    "slui.exe",
    "cmstp.exe"
)
| where InitiatingProcessFileName in~ ("powershell.exe", "pwsh.exe", "cmd.exe", "wscript.exe")
| join kind=leftouter (
    DeviceRegistryEvents
    | where TimeGenerated > ago(1h)
    | where RegistryKey has_any (
        "ms-settings\\shell\\open\\command",
        "mscfile\\shell\\open\\command"
    )
    | summarize RegistryMods = count() by DeviceName, bin(TimeGenerated, 1m)
) on DeviceName
| where isnotnull(RegistryMods) or InitiatingProcessCommandLine has_any ("bypass", "hidden", "-w hidden")
| project TimeGenerated, DeviceName, AccountName, FileName,
    ProcessCommandLine, InitiatingProcessFileName, InitiatingProcessCommandLine
| extend
    Severity = "High",
    ThreatType = "Potential UAC Bypass via Auto-Elevate Binary",
    MitreAttack = "T1548.002",
    Confidence = "Medium"


// ============================================================================
// QUERY 10: Behavioral Correlation - SafePay Attack Chain
// Purpose: Correlate multiple indicators for high-confidence detection
// Confidence: HIGH
// False Positives: Very rare (requires multiple correlated indicators)
// ============================================================================

let UACBypassEvents = DeviceProcessEvents
| where TimeGenerated > ago(2h)
| where ProcessCommandLine has_any ("CMSTPLUA", "3E5FC7F9-9A51-4367-9063-A120244FBEC7", "GetTypeFromCLSID")
| summarize UACAttempts = count() by DeviceName, bin(TimeGenerated, 10m)
| project DeviceName, UACTime = TimeGenerated, UACAttempts;
let RegistryPersistence = DeviceRegistryEvents
| where TimeGenerated > ago(2h)
| where RegistryKey has "CurrentVersion\\Run"
| where ActionType == "RegistryValueSet"
| summarize RegMods = count() by DeviceName, bin(TimeGenerated, 10m)
| project DeviceName, RegTime = TimeGenerated, RegMods;
let ScriptExecution = DeviceProcessEvents
| where TimeGenerated > ago(2h)
| where FileName in~ ("powershell.exe", "pwsh.exe")
| where ProcessCommandLine has_all ("-ExecutionPolicy", "Bypass")
| summarize ScriptRuns = count() by DeviceName, bin(TimeGenerated, 10m)
| project DeviceName, ScriptTime = TimeGenerated, ScriptRuns;
let DefenderTampering = DeviceProcessEvents
| where TimeGenerated > ago(2h)
| where ProcessCommandLine has_any ("Set-MpPreference", "Windows Defender", "DisableRealtimeMonitoring")
| summarize TamperAttempts = count() by DeviceName, bin(TimeGenerated, 10m)
| project DeviceName, TamperTime = TimeGenerated, TamperAttempts;
// Correlate all indicators
UACBypassEvents
| join kind=leftouter (RegistryPersistence) on DeviceName
| join kind=leftouter (ScriptExecution) on DeviceName
| join kind=leftouter (DefenderTampering) on DeviceName
| extend IndicatorCount = (
    iff(UACAttempts > 0, 1, 0) +
    iff(isnotnull(RegMods), 1, 0) +
    iff(isnotnull(ScriptRuns), 1, 0) +
    iff(isnotnull(TamperAttempts), 1, 0)
)
| where IndicatorCount >= 2
| project DeviceName, TimeGenerated = UACTime, UACAttempts, RegMods, ScriptRuns, TamperAttempts,
    IndicatorCount,
    Severity = case(IndicatorCount >= 3, "Critical", "High"),
    ThreatType = "SafePay-Style Attack Chain Detected",
    MitreAttack = "T1548.002, T1547.001, T1562.001",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Confidence = "High"


// ============================================================================
// QUERY 11: Defender Real-Time Protection Status Changes
// Purpose: Detect changes to Defender real-time protection status
// Confidence: HIGH
// False Positives: Legitimate admin actions, policy changes
// MITRE ATT&CK: T1562.001
// ============================================================================

DeviceRegistryEvents
| where TimeGenerated > ago(1h)
| where RegistryKey has "Windows Defender\\Real-Time Protection"
| where ActionType == "RegistryValueSet"
| where RegistryValueName in~ (
    "DisableRealtimeMonitoring",
    "DisableBehaviorMonitoring",
    "DisableOnAccessProtection",
    "DisableScanOnRealtimeEnable"
)
| where RegistryValueData == "1"
| project TimeGenerated, DeviceName, RegistryKey, RegistryValueName,
    RegistryValueData, InitiatingProcessFileName, InitiatingProcessCommandLine,
    AccountName
| extend
    Severity = "Critical",
    ThreatType = "Windows Defender Real-Time Protection Disabled",
    MitreAttack = "T1562.001",
    Confidence = "High"


// ============================================================================
// QUERY 12: PowerShell Add-Type with P/Invoke Definitions
// Purpose: Detect PowerShell loading native Windows APIs (common in UAC bypass)
// Confidence: MEDIUM
// False Positives: Legitimate automation scripts, system management tools
// MITRE ATT&CK: T1106
// ============================================================================

DeviceProcessEvents
| where TimeGenerated > ago(1h)
| where FileName in~ ("powershell.exe", "pwsh.exe")
| where ProcessCommandLine has "Add-Type"
| where ProcessCommandLine has_any (
    "DllImport",
    "user32.dll",
    "kernel32.dll",
    "advapi32.dll",
    "FindWindow",
    "ShowWindow",
    "SendMessage",
    "CreateProcess",
    "OpenProcess"
)
| project TimeGenerated, DeviceName, AccountName, FileName,
    ProcessCommandLine, InitiatingProcessFileName, FolderPath
| extend
    Severity = case(
        ProcessCommandLine has_all ("FindWindow", "ShowWindow"), "High",
        ProcessCommandLine has "CreateProcess", "High",
        "Medium"
    ),
    ThreatType = "PowerShell Native API Usage - Potential Evasion",
    MitreAttack = "T1106",
    Confidence = "Medium"


// ============================================================================
// HUNTING QUERY: F0RT1KA Test Execution Detection
// Purpose: Detect any F0RT1KA security test execution
// Confidence: HIGH
// False Positives: Authorized security testing (coordinate with security team)
// ============================================================================

DeviceFileEvents
| where TimeGenerated > ago(24h)
| where FolderPath startswith "c:\\F0" or FolderPath startswith "C:\\F0"
| summarize
    FileCount = count(),
    FileTypes = make_set(extract(@"\.(\w+)$", 1, FileName)),
    Files = make_set(FileName),
    FirstSeen = min(TimeGenerated),
    LastSeen = max(TimeGenerated)
    by DeviceName, AccountName
| extend
    Severity = "Informational",
    ThreatType = "F0RT1KA Security Test Activity",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11",
    Note = "Verify this is authorized testing activity"
| project DeviceName, AccountName, FileCount, FileTypes, Files, FirstSeen, LastSeen,
    Severity, ThreatType, TestID, Note


// ============================================================================
// ALERT SUMMARY DASHBOARD QUERY
// Purpose: Aggregate all SafePay-related detections for SOC dashboard
// ============================================================================

let DetectionEvents = union
    (DeviceProcessEvents
    | where TimeGenerated > ago(24h)
    | where ProcessCommandLine has_any ("CMSTPLUA", "3E5FC7F9-9A51-4367-9063-A120244FBEC7")
    | summarize EventCount = count() by Computer = DeviceName, bin(TimeGenerated, 1h)
    | extend DetectionType = "UAC Bypass Attempt"
    ),
    (DeviceRegistryEvents
    | where TimeGenerated > ago(24h)
    | where RegistryKey has "CurrentVersion\\Run"
    | where RegistryValueName =~ "SafePayService" or RegistryValueData has "6F22-C16F"
    | summarize EventCount = count() by Computer = DeviceName, bin(TimeGenerated, 1h)
    | extend DetectionType = "SafePay Persistence"
    ),
    (DeviceProcessEvents
    | where TimeGenerated > ago(24h)
    | where ProcessCommandLine has_any ("Set-MpPreference", "DisableRealtimeMonitoring")
    | summarize EventCount = count() by Computer = DeviceName, bin(TimeGenerated, 1h)
    | extend DetectionType = "Defender Tampering"
    ),
    (DeviceFileEvents
    | where TimeGenerated > ago(24h)
    | where FolderPath startswith "c:\\F0"
    | where FileName endswith ".ps1"
    | summarize EventCount = count() by Computer = DeviceName, bin(TimeGenerated, 1h)
    | extend DetectionType = "Script Drop to c:\\F0"
    );
DetectionEvents
| summarize TotalEvents = sum(EventCount), DetectionTypes = make_set(DetectionType)
    by Computer, bin(TimeGenerated, 1h)
| order by TotalEvents desc
| project TimeGenerated, Computer, TotalEvents, DetectionTypes,
    Severity = case(TotalEvents >= 5, "Critical", TotalEvents >= 3, "High", "Medium"),
    ThreatType = "SafePay Attack Activity Cluster",
    MitreAttack = "T1548.002, T1547.001, T1562.001",
    TestID = "2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11"
