# tests_source/intel-driven/2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11/2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11_dr_rules.yaml
# ============================================================================
# LimaCharlie Detection & Response Rules
# ============================================================================
# Test ID: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
# Test Name: SafePay UAC Bypass & Defense Evasion
# MITRE ATT&CK: T1548.002, T1562.001, T1547.001
# Created: 2025-12-07
# Author: F0RT1KA Defense Guidance Builder
# ============================================================================
#
# FILE FORMAT: For use with `limacharlie dr add` or `config push`
#
# Deployment:
#   limacharlie dr add --rule-file 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11_dr_rules.yaml
#
# Or with limacharlie-cli config push:
#   limacharlie config push dr --rule-file 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11_dr_rules.yaml
#
# ============================================================================

rules:
  # ==========================================================================
  # Rule 1: PowerShell Execution Policy Bypass with SafePay Indicators
  # Purpose: Detect PowerShell script execution with bypass containing SafePay indicators
  # Confidence: HIGH
  # MITRE ATT&CK: T1059.001
  # ==========================================================================
  safepay-powershell-bypass:
    detect:
      target: edr
      event: NEW_PROCESS
      op: and
      rules:
        - op: is
          path: event/FILE_PATH
          value: powershell.exe
          case sensitive: false
        - op: contains
          path: event/COMMAND_LINE
          value: "-ExecutionPolicy"
          case sensitive: false
        - op: contains
          path: event/COMMAND_LINE
          value: "Bypass"
          case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/COMMAND_LINE
              value: "safepay"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "uac_bypass"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "cmstplua"
              case sensitive: false

    respond:
      - action: report
        name: SafePay PowerShell Bypass Execution
        metadata:
          severity: high
          mitre_attack: T1059.001
          description: PowerShell script executed with execution policy bypass containing SafePay-related indicators


  # ==========================================================================
  # Rule 2: CMSTPLUA COM Object UAC Bypass Detection
  # Purpose: Detect UAC bypass via CMSTPLUA COM interface
  # Confidence: HIGH
  # MITRE ATT&CK: T1548.002
  # ==========================================================================
  safepay-cmstplua-uac-bypass:
    detect:
      target: edr
      event: NEW_PROCESS
      op: and
      rules:
        - op: or
          rules:
            - op: is
              path: event/FILE_PATH
              value: powershell.exe
              case sensitive: false
            - op: is
              path: event/FILE_PATH
              value: pwsh.exe
              case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/COMMAND_LINE
              value: "3E5FC7F9-9A51-4367-9063-A120244FBEC7"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "CMSTPLUA"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "GetTypeFromCLSID"
              case sensitive: false

    respond:
      - action: report
        name: CMSTPLUA UAC Bypass Attempt
        metadata:
          severity: critical
          mitre_attack: T1548.002
          test_id: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
          description: Detected attempt to bypass UAC using CMSTPLUA COM object


  # ==========================================================================
  # Rule 3: SafePay Registry Persistence - Specific IoC
  # Purpose: Detect SafePay-specific registry persistence indicator
  # Confidence: HIGH
  # MITRE ATT&CK: T1547.001
  # ==========================================================================
  safepay-registry-persistence-ioc:
    detect:
      target: edr
      event: REGISTRY_VALUE_SET
      op: and
      rules:
        - op: contains
          path: event/REGISTRY_KEY
          value: "CurrentVersion\\Run"
          case sensitive: false
        - op: or
          rules:
            - op: is
              path: event/VALUE_NAME
              value: SafePayService
              case sensitive: false
            - op: contains
              path: event/VALUE_DATA
              value: "6F22-C16F-0C71-688A"
              case sensitive: false

    respond:
      - action: report
        name: SafePay Malware Persistence Detected
        metadata:
          severity: critical
          mitre_attack: T1547.001
          test_id: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
          description: SafePay-specific registry persistence indicator detected


  # ==========================================================================
  # Rule 4: Generic Run Key Persistence from Suspicious Process
  # Purpose: Detect Run key modifications from processes in user-writable locations
  # Confidence: MEDIUM
  # MITRE ATT&CK: T1547.001
  # ==========================================================================
  safepay-run-key-suspicious-process:
    detect:
      target: edr
      event: REGISTRY_VALUE_SET
      op: and
      rules:
        - op: contains
          path: event/REGISTRY_KEY
          value: "CurrentVersion\\Run"
          case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/PROCESS_PATH
              value: "\\Temp\\"
              case sensitive: false
            - op: contains
              path: event/PROCESS_PATH
              value: "\\AppData\\Local\\Temp"
              case sensitive: false
            - op: contains
              path: event/PROCESS_PATH
              value: "\\Users\\Public"
              case sensitive: false
            - op: contains
              path: event/PROCESS_PATH
              value: "\\Downloads\\"
              case sensitive: false

    respond:
      - action: report
        name: Suspicious Run Key Persistence Attempt
        metadata:
          severity: high
          mitre_attack: T1547.001
          description: Registry Run key modified by process in user-writable location (common malware staging area)


  # ==========================================================================
  # Rule 5: Windows Defender Tampering via PowerShell
  # Purpose: Detect attempts to disable Windows Defender protection
  # Confidence: HIGH
  # MITRE ATT&CK: T1562.001
  # ==========================================================================
  safepay-defender-tampering:
    detect:
      target: edr
      event: NEW_PROCESS
      op: and
      rules:
        - op: or
          rules:
            - op: is
              path: event/FILE_PATH
              value: powershell.exe
              case sensitive: false
            - op: is
              path: event/FILE_PATH
              value: pwsh.exe
              case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/COMMAND_LINE
              value: "Set-MpPreference"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "DisableRealtimeMonitoring"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "DisableBehaviorMonitoring"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "DisableIOAVProtection"
              case sensitive: false

    respond:
      - action: report
        name: Windows Defender Tampering Attempt
        metadata:
          severity: critical
          mitre_attack: T1562.001
          test_id: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
          description: Attempt to disable Windows Defender protection via PowerShell


  # ==========================================================================
  # Rule 6: Suspicious Script File Drop with SafePay Indicators
  # Purpose: Detect script creation with SafePay-related naming patterns
  # Confidence: HIGH
  # ==========================================================================
  safepay-suspicious-script-drop:
    detect:
      target: edr
      event: FILE_CREATE
      op: and
      rules:
        - op: or
          rules:
            - op: ends with
              path: event/FILE_PATH
              value: ".ps1"
              case sensitive: false
            - op: ends with
              path: event/FILE_PATH
              value: ".bat"
              case sensitive: false
            - op: ends with
              path: event/FILE_PATH
              value: ".cmd"
              case sensitive: false
            - op: ends with
              path: event/FILE_PATH
              value: ".vbs"
              case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/FILE_PATH
              value: "safepay"
              case sensitive: false
            - op: contains
              path: event/FILE_PATH
              value: "uac_bypass"
              case sensitive: false
            - op: contains
              path: event/FILE_PATH
              value: "persistence"
              case sensitive: false

    respond:
      - action: report
        name: Suspicious SafePay-Related Script Dropped
        metadata:
          severity: high
          mitre_attack: T1105
          description: Script file created with SafePay-related naming pattern


  # ==========================================================================
  # Rule 7: GUI Automation for Security Tool Interaction
  # Purpose: Detect GUI automation targeting Windows Security
  # Confidence: MEDIUM
  # MITRE ATT&CK: T1562.001
  # ==========================================================================
  safepay-gui-automation-defender:
    detect:
      target: edr
      event: NEW_PROCESS
      op: and
      rules:
        - op: or
          rules:
            - op: is
              path: event/FILE_PATH
              value: powershell.exe
              case sensitive: false
            - op: is
              path: event/FILE_PATH
              value: pwsh.exe
              case sensitive: false
        - op: contains
          path: event/COMMAND_LINE
          value: "Add-Type"
          case sensitive: false
        - op: or
          rules:
            - op: contains
              path: event/COMMAND_LINE
              value: "FindWindow"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "ShowWindow"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "user32.dll"
              case sensitive: false
            - op: contains
              path: event/COMMAND_LINE
              value: "UIAutomation"
              case sensitive: false

    respond:
      - action: report
        name: GUI Automation Targeting Security Tools
        metadata:
          severity: high
          mitre_attack: T1562.001
          test_id: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
          description: PowerShell script using GUI automation APIs potentially targeting security tools


  # ==========================================================================
  # Rule 8: Defender Real-Time Protection Registry Modification
  # Purpose: Detect registry changes to Defender real-time protection settings
  # Confidence: HIGH
  # MITRE ATT&CK: T1562.001
  # ==========================================================================
  safepay-defender-registry-tampering:
    detect:
      target: edr
      event: REGISTRY_VALUE_SET
      op: and
      rules:
        - op: contains
          path: event/REGISTRY_KEY
          value: "Windows Defender\\Real-Time Protection"
          case sensitive: false
        - op: or
          rules:
            - op: is
              path: event/VALUE_NAME
              value: DisableRealtimeMonitoring
              case sensitive: false
            - op: is
              path: event/VALUE_NAME
              value: DisableBehaviorMonitoring
              case sensitive: false
            - op: is
              path: event/VALUE_NAME
              value: DisableOnAccessProtection
              case sensitive: false
        - op: is
          path: event/VALUE_DATA
          value: "1"

    respond:
      - action: report
        name: Windows Defender Protection Disabled via Registry
        metadata:
          severity: critical
          mitre_attack: T1562.001
          test_id: 2cf59d3e-ae82-48bb-9779-4a5ba5bd9c11
          description: Registry modification detected that disables Windows Defender real-time protection
